# layout/DEBIAN/postinst
launchctl load /fs/jb/Library/LaunchDaemons/com.nito.tssagent.plist
